# standalone/service/pkg-index/run
#!/bin/sh
. $HOME/setup-env
cd $HOME/pkg-index/official
exec racket configs/live.rkt 2>&1
