# service/pkgserver/run
#!/bin/sh
make -C /usr/local/racket-pkg-website keys
export CONFIG=docker
exec /sbin/setuser pkgserver /bin/sh -c 'cd /usr/local/racket-pkg-website; ./run'
