# service/pkgserver/run
#!/bin/sh
cd /usr/local/racket-pkg-website
make keys
export CONFIG=docker
exec /sbin/setuser pkgserver ./run
