# standalone/service/racket-pkg-website/run
#!/bin/sh
. $HOME/setup-env
cd $HOME/racket-pkg-website
export CONFIG=live
exec ./run
