# service/pkg-index/run
#!/bin/sh
cd /usr/local/pkg-index/official
chown -R pkgserver:pkgserver /var/lib/pkgserver
exec /sbin/setuser pkgserver racket configs/docker.rkt
