# service/log-forwarder/run
#!/bin/sh
sv start pkg-index || exit 1
sv start racket-pkg-website || exit 1
exec tail -F -n 0 /var/log/pkgserver/*/current
